# audio/alsa-utils-1.0.29/alsactl/90-alsa-restore.rules
ACTION=="add", SUBSYSTEM=="sound", KERNEL=="controlC*", KERNELS!="card*", TEST=="/home/tao/work/sama5d3/pos/audio/install/sbin", TEST=="/home/tao/work/sama5d3/pos/audio/install/share/alsa", GOTO="alsa_restore_go"
GOTO="alsa_restore_end"

LABEL="alsa_restore_go"
TEST!="/etc/alsa/state-daemon.conf", RUN+="/home/tao/work/sama5d3/pos/audio/install/sbin/alsactl restore $attr{device/number}"
TEST=="/etc/alsa/state-daemon.conf", RUN+="/home/tao/work/sama5d3/pos/audio/install/sbin/alsactl nrestore $attr{device/number}"

LABEL="alsa_restore_end"
